# onsite_revenue_management/templates/default/revenue_management.crontab.erb
# Note: OpsWorks Instances are UTC
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
RAILS_ENV=<% @environment %>
*/2 * * * * cd <%= @rails_path %> && bundle exec rake rx:crontest
11 8 * * * cd  <%= @rails_path %> && bundle exec rake rx:publish
